# docker/supervisor/supervisord.conf
# =============================================================================
# Supervisor Configuration for Laravel Genealogy Application
# =============================================================================
# Manages Nginx, PHP-FPM, and Laravel queue workers

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=true
silent=false
minfds=1024
minprocs=200
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# =============================================================================
# Nginx Web Server
# =============================================================================
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/log/supervisor/nginx_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/nginx_stdout.log
stdout_logfile_maxbytes=10MB

# =============================================================================
# PHP-FPM Application Server
# =============================================================================
[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize --fpm-config /usr/local/etc/php-fpm.conf
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/log/supervisor/php-fpm_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/php-fpm_stdout.log
stdout_logfile_maxbytes=10MB

# =============================================================================
# Laravel Queue Workers for Genealogy Processing
# =============================================================================
[program:laravel-worker-default]
command=/usr/local/bin/php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600 --queue=default
process_name=%(program_name)s_%(process_num)02d
numprocs=2
directory=/var/www/html
autostart=true
autorestart=true
startretries=5
startsecs=1
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-worker-default.log
stdout_logfile_maxbytes=10MB
stopwaitsecs=3600
killasgroup=true
stopasgroup=true

# Media processing queue for genealogy photos
[program:laravel-worker-media]
command=/usr/local/bin/php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=1800 --queue=media-conversions
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/var/www/html
autostart=true
autorestart=true
startretries=5
startsecs=1
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-worker-media.log
stdout_logfile_maxbytes=10MB
stopwaitsecs=1800
killasgroup=true
stopasgroup=true

# Backup processing queue
[program:laravel-worker-backup]
command=/usr/local/bin/php /var/www/html/artisan queue:work --sleep=10 --tries=1 --max-time=7200 --queue=backups
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/var/www/html
autostart=true
autorestart=true
startretries=5
startsecs=1
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-worker-backup.log
stdout_logfile_maxbytes=10MB
stopwaitsecs=7200
killasgroup=true
stopasgroup=true

# =============================================================================
# Laravel Task Scheduler
# =============================================================================
[program:laravel-scheduler]
command=/bin/bash -c "while [ true ]; do (php /var/www/html/artisan schedule:run --verbose --no-interaction &); sleep 60; done"
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/var/www/html
autostart=true
autorestart=true
startretries=5
startsecs=1
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-scheduler.log
stdout_logfile_maxbytes=10MB
killasgroup=true
stopasgroup=true

# =============================================================================
# Laravel Horizon (if Redis queue is used)
# =============================================================================
[program:laravel-horizon]
command=/usr/local/bin/php /var/www/html/artisan horizon
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/var/www/html
autostart=false
autorestart=true
startretries=5
startsecs=1
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-horizon.log
stdout_logfile_maxbytes=10MB
stopwaitsecs=3600
killasgroup=true
stopasgroup=true

# =============================================================================
# Log Management
# =============================================================================
[program:log-rotator]
command=/bin/bash -c "while [ true ]; do find /var/www/html/storage/logs -name '*.log' -size +100M -exec truncate -s 50M {} \;; sleep 3600; done"
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
startretries=5
startsecs=1
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/log-rotator.log
stdout_logfile_maxbytes=1MB

# =============================================================================
# Groups for easy management
# =============================================================================
[group:web]
programs=nginx,php-fpm

[group:laravel]
programs=laravel-worker-default,laravel-worker-media,laravel-worker-backup,laravel-scheduler

[group:maintenance]
programs=log-rotator